# workspace/crontab/cantor.job
* * * * * cantor /usr/bin/php /srv/services/time-words-server/artisan schedule:run >> /dev/null 2>&1
