# workspace/crontab/cantor.job.example
* * * * * cantor /usr/bin/php /srv/services/artisan schedule:run >> /dev/null 2>&1
